# scripts/files-to-agent-update-host.service
[Unit]
Description=files-to-agent update watcher (host side)
After=docker.service
Requires=docker.service

[Service]
Type=simple
# Adjust paths to your install location.
Environment=COMPOSE_DIR=/srv/files-to-agent
Environment=FLAG_DIR=/srv/files-to-agent/update-flag
Environment=POLL_SECONDS=5
ExecStart=/usr/local/bin/files-to-agent-update-host
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
